# services/restartnfc.timer
[Unit]
Description=Restart nfc reader

[Timer]
OnActiveSec=10s
OnBootSec=60s
OnUnitActiveSec=24h 

[Install]
WantedBy=timers.target

